# com.kijito.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kijito.monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>-u</string>
    <string>/Users/jason/Code/Kijito.ai/kijito_monitor/monitor/kijito_monitor.py</string>
    <string>--state-file</string>
    <string>/Users/jason/.cache/kijito-monitor/hive.json</string>
    <string>--events-file</string>
    <string>/Users/jason/.cache/kijito-monitor/events.ndjson</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/jason/Code/Kijito.ai/kijito_monitor/monitor</string>
  <key>StandardOutPath</key>
  <string>/Users/jason/.cache/kijito-monitor/monitor.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/jason/.cache/kijito-monitor/monitor.err</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
